# supervisor_inis/apostify.ini
[program:apostify]

command=/home/python/bin/uwsgi /home/apostify/uwsgi.ini
environment=OPENSSL_CONF=/home/apostify/openssl.cnf
user=root
stopsignal=QUIT
autostart=true
autorestart=true
stdout_logfile=/home/apostify/logs/supervisor.log
stderr_logfile=/home/apostify/logs/supervisor_err.log
stopasgroup=true
killasgroup=true