# supervisor_inis/celery_worker.ini
[program:celery-worker]

command=celery -A core worker -l INFO --concurrency=4
directory=/home/apostify/core
numprocs=1
stdout_logfile=/home/apostify/logs/celery_worker.log
stderr_logfile=/home/apostify/logs/celery_worker_error.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
priority=998
