# supervisor_inis/celery_beat.ini
[program:celery-beat]

command=celery -A core beat --scheduler django_celery_beat.schedulers:DatabaseScheduler -l info
directory=/home/apostify/core
numprocs=1
stdout_logfile=/home/apostify/logs/celery_beat.log
stderr_logfile=/home/apostify/logs/celery_beat_error.log
autostart=true
autorestart=true
startsecs=10
priority=999